# airootfs/etc/skel/.config/autostart/xcompmgr.desktop
[Desktop Entry]
Exec=/usr/bin/xcompmgr
Name=xcompmgr
Type=Application
Version=1.0
